# docker/prod/cron/crontab
*    *   * * * www-data cd /api &&  timeout 300 php app/console rec:crypto:check --env=prod
*    *   * * * www-data cd /api &&  timeout 300 php app/console rec:fiat:check --env=prod
*    *   * * * www-data cd /api &&  timeout 300 php app/console rec:check:offers --env=prod
*    *   * * * www-data find /api/web/static -name *.tmp.* -type f -mmin +10 -delete
